# debian/postinst
#!/bin/sh

# override tool behaviour through distro-defaults
FLL_DISTRO_MODE="installed"
if [ -r /etc/default/distro ]; then
	source /etc/default/distro
fi

if [ ! "$FLL_DISTRO_MODE" = "live" ]; then
	rm -f /etc/X11/Xsession.d/45xsession
	exit 0
fi

case "$1" in
    configure)
		# This must be in /etc/inittab, because the script does not return.
		#update-rc.d xsession start 99 5 . stop 00 0 1 2 3 4 6 .
		update-rc.d xsession stop 10 0 1 2 3 4 6 .
		grep -q "x5:5:wait:/etc/init.d/xsession start" /etc/inittab && exit 0
		cat >>/etc/inittab <<EOF

#-- xsession-sidux start
# Run X Window session from CDROM in runlevel 5
w5:5:wait:/bin/sleep 2
x5:5:wait:/etc/init.d/xsession start
#-- xsession-sidux end
EOF
		;;

	abort-upgrade|abort-remove|abort-deconfigure)
		 ;;

	*)
		echo "postinst called with unknown argument \`$1'" >&2
		exit 0
		;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

