# debian/postrm
#!/bin/sh
case "$1" in
	purge|remove|upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
		# This must be in /etc/inittab because xsession doesn't return
		update-rc.d -f xsession remove 2>/dev/null

		perl -pi -e 's/^# Run X Window session from CDROM in runlevel 5\n//g;
	             s|^w5:5:wait:/bin/sleep 2\n||g;
	             s|^x5:5:wait:/etc/init.d/xsession start\n||g;' /etc/inittab
	        ;;
	*)
		echo "postrm called with unknown argument \`$1'" >&2
		exit 0
		;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

